# specific/arch/base/.config/systemd/user/claude-remote.service
[Unit]
Description=Claude code remote control
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/khaosdoctor/claude-home
ExecStart=bash -lc 'echo y | claude remote-control --spawn same-dir --name "Persistent %H session"'
StandardOutput=null
Restart=always
RestartSec=10
TimeoutStartSec=60
Nice=10

[Install]
WantedBy=default.target
