# specific/arch/base/.config/systemd/user/claude-remote.service
[Unit]
Description=Claude code remote control
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/script -qfc "/usr/bin/claude --remote-control --no-chrome --name 'Persistent %H session'" /dev/null
Restart=always
RestartSec=10
TimeoutStartSec=60
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=default.target
