# specific/arch/base/.config/systemd/user/clearcache.timer
[Unit]
Description=Run cache cleanup weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
