# specific/arch/base/.config/systemd/user/elephant.service
[Unit]
Description=Run elephant on Hyprland startup
PartOf=graphical-session.target
After=graphical-session.target
StartLimitIntervalSec=30
StartLimitBurst=20

[Service]
Type=simple
ExecStart=/usr/bin/elephant
Restart=on-failure
RestartSec=0.5

[Install]
WantedBy=graphical-session.target
