# specific/arch/base/.config/systemd/user/graphical-session-apps.target
[Unit]
Description=Post-session application services
BindsTo=graphical-session.target
After=graphical-session.target
