# specific/arch/base/.config/systemd/user/obsidian-agent@.service
[Unit]
Description=Obsidian second-brain agent (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
TimeoutStartSec=1800
ExecStart=%h/.claude/bin/obsidian-agent %i
