# specific/arch/base/.config/systemd/user/rclone-b2-encrypted-backups.timer
[Unit]
Description=Weekly encrypted sync to Backblaze B2

[Timer]
OnCalendar=Wed 19:00
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
