# specific/arch/base/.config/systemd/user/wallpaper-changer.timer
[Unit]
Description=Run wallpaper change every hour

[Timer]
OnUnitActiveSec=1h
Unit=wallpaper-changer.service

[Install]
WantedBy=timers.target
